# devconfigs/rpi4/files/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.ipaddr='172.17.17.2'
set network.One=interface
set network.One.proto='dhcp'
set network.One.device='1-USB'
set network.One.usename='1-USB'
set network.One.metric='10'
set network.Two=interface
set network.Two.proto='dhcp'
set network.Two.device='2-USB'
set network.Two.usename='2-USB'
set network.Two.metric='20'
set network.Three=interface
set network.Three.proto='dhcp'
set network.Three.device='3-USB'
set network.Three.usename='3-USB'
set network.Three.metric='30'
set network.Four=interface
set network.Four.proto='dhcp'
set network.Four.device='4-USB'
set network.Four.usename='4-USB'
set network.Four.metric='40'
commit network
set system.@system[-1].hostname='SmoothWAN'
add system led
set system.@led[0]=led
set system.@led[0].name='System Load Indicator'
set system.@led[0].sysfs='led0'
set system.@led[0].trigger='heartbeat'
add system led
set system.@led[1].name='Boot Done'
set system.@led[1].sysfs='led1'
set system.@led[1].trigger='none'
commit system
set dhcp.cfg01411c.allservers='1'
add_list dhcp.cfg01411c.server='8.8.8.8'
add_list dhcp.cfg01411c.server='8.8.4.4'
set dhcp.cfg01411c.rebind_protection='0'
commit dhcp
del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='Speedify'
add_list firewall.cfg03dc81.network='One'
add_list firewall.cfg03dc81.network='Two'
add_list firewall.cfg03dc81.network='Three'
add_list firewall.cfg03dc81.network='Four'
commit firewall
set network.Speedify=interface
set network.Speedify.proto='none'
set network.Speedify.device='connectify0'
commit network
set nlbwmon.cfg015d22.database_directory='/usr/share/nlbwmon'
set nlbwmon.cfg015d22.commit_interval='10m'
commit nlbwmon

EOI

