# devconfigs/rpi4/files/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.ipaddr='172.17.17.2'
set network.One=interface
set network.One.proto='dhcp'
set network.One.device='1-USB'
set network.One.usename='1-USB'
set network.One.metric='10'
set network.Two=interface
set network.Two.proto='dhcp'
set network.Two.device='2-USB'
set network.Two.usename='2-USB'
set network.Two.metric='20'
set network.Three=interface
set network.Three.proto='dhcp'
set network.Three.device='3-USB'
set network.Three.usename='3-USB'
set network.Three.metric='30'
set network.Four=interface
set network.Four.proto='dhcp'
set network.Four.device='4-USB'
set network.Four.usename='4-USB'
set network.Four.metric='40'
commit network
set system.@system[-1].hostname='SmoothWAN'
add system led
set system.@led[0]=led
set system.@led[0].name='System Load Indicator'
set system.@led[0].sysfs='led0'
set system.@led[0].trigger='heartbeat'
add system led
set system.@led[1].name='Boot Done'
set system.@led[1].sysfs='led1'
set system.@led[1].trigger='none'
commit system
set dhcp.cfg01411c.allservers='1'
set dhcp.cfg01411c.rebind_protection='0'
set dhcp.lan.force='1'
del dhcp.lan.ra
del dhcp.lan.ra_default
del dhcp.lan.ra_slaac
del dhcp.lan.dhcpv6
del dhcp.lan.ra_flags
add_list dhcp.lan.ra_flags='none'
set dhcp.@dnsmasq[0].dhcpscript='/usr/sbin/iptmon'
commit dhcp
del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='Speedify'
add_list firewall.cfg03dc81.network='One'
add_list firewall.cfg03dc81.network='Two'
add_list firewall.cfg03dc81.network='Three'
add_list firewall.cfg03dc81.network='Four'
commit firewall
set network.Speedify=interface
set network.Speedify.proto='none'
set network.Speedify.device='connectify0'
del network.cfg030f15.ports
add_list network.cfg030f15.ports='eth0'
add_list network.cfg030f15.ports='usb0'
commit network
set system.ntp.use_dhcp='0'
del system.ntp.server
add_list system.ntp.server='162.159.200.123'
commit system
set luci.apply.holdoff='1'
commit luci
add_list vpn-policy-routing.config.supported_interface='One'
add_list vpn-policy-routing.config.supported_interface='Two'
add_list vpn-policy-routing.config.supported_interface='Three'
add_list vpn-policy-routing.config.supported_interface='Four'
add_list vpn-policy-routing.config.supported_interface='Speedify'
commit vpn-policy-routing
set luci_statistics.collectd_sensors.enable='1'
set luci_statistics.collectd_sensors.IgnoreSelected='1'
del luci_statistics.cfg1b3a2e
set luci_statistics.collectd_iptables.enable='1'
set luci_statistics.cfg1a3a2e.name='LAN'
set luci_statistics.cfg1a3a2e.chain='zone_lan_postrouting'
set luci_statistics.collectd_ping.enable='1'
del luci_statistics.collectd_ping.Hosts
add_list luci_statistics.collectd_ping.Hosts='10.202.0.1'
set luci_statistics.collectd_ping.AddressFamily='any'
del luci_statistics.collectd_uptime.enable
del luci_statistics.collectd_network.enable
set luci_statistics.collectd.Include='/etc/collectd/conf.d/*.conf'
set luci_statistics.collectd_uptime.enable='1'
del luci_statistics.collectd_interface.IgnoreSelected
del luci_statistics.collectd_interface.Interfaces
add_list luci_statistics.collectd_interface.Interfaces='3-USB'
add_list luci_statistics.collectd_interface.Interfaces='br-lan'
add_list luci_statistics.collectd_interface.Interfaces='connectify0'
add_list luci_statistics.collectd_interface.Interfaces='1-USB'
add_list luci_statistics.collectd_interface.Interfaces='2-USB'
add_list luci_statistics.collectd_interface.Interfaces='4-USB'
commit luci_statistics

EOI

#Replace diagnostics page
cp /etc/sysupgdiff/diagnostics.js /www/luci-static/resources/view/network/diagnostics.js