# devconfigs/flint/files/etc/uci-defaults/99-custom
uci -q batch << EOI
del network.wan
del network.wan6

set network.lan.ipaddr='172.17.17.2'

set network.PortWAN=interface
set network.PortWAN.proto='dhcp'
set network.PortWAN.device='eth0'
set network.PortWAN.metric='10'

set network.PortOne=interface
set network.PortOne.proto='dhcp'
set network.PortOne.device='eth1'
set network.PortOne.metric='20'

set network.PortTwo=interface
set network.PortTwo.proto='dhcp'
set network.PortTwo.device='eth2'
set network.PortTwo.metric='30'

set network.PortThree=interface
set network.PortThree.proto='dhcp'
set network.PortThree.device='eth3'
set network.PortThree.metric='40'

set network.AutoUSBeth5=interface
set network.AutoUSBeth5.proto='dhcp'
set network.AutoUSBeth5.device='eth5'
set network.AutoUSBeth5.metric='35'

set network.AutoUSBeth6=interface
set network.AutoUSBeth6.proto='dhcp'
set network.AutoUSBeth6.device='eth6'
set network.AutoUSBeth6.metric='36'

set network.AutoUSBeth7=interface
set network.AutoUSBeth7.proto='dhcp'
set network.AutoUSBeth7.device='eth7'
set network.AutoUSBeth7.metric='37'

set network.AutoUSB0=interface
set network.AutoUSB0.proto='dhcp'
set network.AutoUSB0.device='usb0'
set network.AutoUSB0.metric='40'

set network.AutoUSB1=interface
set network.AutoUSB1.proto='dhcp'
set network.AutoUSB1.device='usb1'
set network.AutoUSB1.metric='41'

set network.AutoUSB2=interface
set network.AutoUSB2.proto='dhcp'
set network.AutoUSB2.device='usb2'
set network.AutoUSB2.metric='42'

set network.AutoUSB3=interface
set network.AutoUSB3.proto='dhcp'
set network.AutoUSB3.device='usb3'
set network.AutoUSB3.metric='43'

set network.wwan=interface
set network.wwan.proto='dhcp'
set network.wwan.metric='100'

set network.wwan2=interface
set network.wwan2.proto='dhcp'
set network.wwan2.metric='99'

del network.cfg030f15.ports
add_list network.cfg030f15.ports='eth4'
commit network

set system.@system[-1].hostname='SmoothWAN'
set system.@system[-1].log_size='1024'
del system.led_wan
add system led
set system.@led[0].trigger='default-on'
set system.@led[0].sysfs='white_led'
set system.@led[0].name='Boot Done'
add system led
set system.@led[1].trigger='none'
set system.@led[1].sysfs='blue_led'
set system.@led[1].name='Boot Done'
commit system

add_list dhcp.lan.dhcp_option='6,172.17.17.2'
set dhcp.@dnsmasq[-1].port='0'
commit dhcp

del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='wan'
add_list firewall.cfg03dc81.network='wwan'
add_list firewall.cfg03dc81.network='wwan2'
add_list firewall.cfg03dc81.network='PortWAN'
add_list firewall.cfg03dc81.network='PortOne'
add_list firewall.cfg03dc81.network='PortTwo'
add_list firewall.cfg03dc81.network='PortThree'
add_list firewall.cfg03dc81.network='AutoUSB0'
add_list firewall.cfg03dc81.network='AutoUSB1'
add_list firewall.cfg03dc81.network='AutoUSB2'
add_list firewall.cfg03dc81.network='AutoUSB3'
add_list firewall.cfg03dc81.network='AutoUSBeth5'
add_list firewall.cfg03dc81.network='AutoUSBeth6'
add_list firewall.cfg03dc81.network='AutoUSBeth7'
commit firewall

del network.globals.ula_prefix
commit network

set system.ntp.use_dhcp='0'
del system.ntp.server
add_list system.ntp.server='162.159.200.123'
commit system

add_list vpn-policy-routing.config.supported_interface='PortWAN'
add_list vpn-policy-routing.config.supported_interface='PortOne'
add_list vpn-policy-routing.config.supported_interface='PortTwo'
add_list vpn-policy-routing.config.supported_interface='PortThree'
add_list vpn-policy-routing.config.supported_interface='AutoUSB0'
add_list vpn-policy-routing.config.supported_interface='AutoUSB1'
add_list vpn-policy-routing.config.supported_interface='AutoUSB2'
add_list vpn-policy-routing.config.supported_interface='AutoUSB3'
add_list vpn-policy-routing.config.supported_interface='AutoUSBeth5'
add_list vpn-policy-routing.config.supported_interface='AutoUSBeth6'
add_list vpn-policy-routing.config.supported_interface='AutoUSBeth7'
add_list vpn-policy-routing.config.supported_interface='Speedify'
add_list vpn-policy-routing.config.supported_interface='wan'
commit vpn-policy-routing

set luci_statistics.collectd_sensors.enable='1'
set luci_statistics.collectd_sensors.IgnoreSelected='1'
set luci_statistics.collectd_ping.enable='1'
del luci_statistics.collectd_ping.Hosts
add_list luci_statistics.collectd_ping.Hosts='10.202.0.1'
add_list luci_statistics.collectd_ping.Hosts='10.202.0.5'
add_list luci_statistics.collectd_ping.Hosts='8.8.8.8'
set luci_statistics.collectd_ping.AddressFamily='any'
del luci_statistics.collectd_uptime.enable
del luci_statistics.collectd_network.enable
set luci_statistics.collectd.Include='/etc/collectd/conf.d/*.conf'
set luci_statistics.collectd_uptime.enable='1'
set luci_statistics.collectd_interface.IgnoreSelected='1'
add_list luci_statistics.collectd_interface.Interfaces='br-lan'
add_list luci_statistics.collectd_interface.Interfaces='ip6tnl0'
add_list luci_statistics.collectd_interface.Interfaces='bond0'
add_list luci_statistics.collectd_interface.Interfaces='lo'
add_list luci_statistics.collectd_interface.Interfaces='gretap0'
add_list luci_statistics.collectd_interface.Interfaces='ip6gre0'
add_list luci_statistics.collectd_interface.Interfaces='gre0'
set luci_statistics.collectd_rrdtool.DataDir='/etc/rrd'
commit luci_statistics
add fstab mount # =cfg024d78
set fstab.@mount[-1].enabled='1'
set fstab.@mount[-1].target='/mnt/smoothwan'
set fstab.@mount[-1].uuid='337d0b9f-8716-42b2-b316-312811aef984'
commit fstab
add_list nlbwmon.cfg015d22.local_network='172.17.17.0/24'
set nlbwmon.cfg015d22.database_compress='1'
set nlbwmon.cfg015d22.database_directory='/usr/lib/nlbwmon'
commit nlbwmon
add firewall zone # =cfg11dc81
set firewall.@zone[-1].name='Tailscale'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'
set firewall.@zone[-1].masq='1'
add_list firewall.@zone[-1].device='br-tailscale0'
add firewall forwarding # =cfg12ad58
set firewall.@forwarding[-1].src='Tailscale'
set firewall.@forwarding[-1].dest='lan'
add firewall forwarding # =cfg13ad58
set firewall.@forwarding[-1].src='Tailscale'
set firewall.@forwarding[-1].dest='wan'
add firewall forwarding # =cfg14ad58
set firewall.@forwarding[-1].src='lan'
set firewall.@forwarding[-1].dest='Tailscale'
set firewall.@zone[1].masq6="1"
commit firewall

EOI


#Replace diagnostics page
cp /etc/sysupgdiff/diagnostics.js /www/luci-static/resources/view/network/diagnostics.js

#Hints
sed -i "s+Choose the network(s) you want to attach to this wireless interface or fill out the <em>custom</em> field to define a new network.+Choose the network(s) you want to attach to this wireless interface or fill out the <em>custom</em> field to define a new network.<br><br><b>SmoothWAN Note!<br><br>If this is a Wireless WAN, head to Network->Interfaces and edit the automatically created interface to set the metric to a unique number (e.g 10 or more).</b>+g" /www/luci-static/resources/view/network/wireless.js
sed -i "s+_('Bring up on boot')+_('Bring up on boot'), _('<br><br><b>SmoothWAN Note!<br><br>If this is a WAN connection, head to Advanced Settings tab above and set metric to a unique number (e.g 10 or more), then to Firewall and set zone to red.')+" /www/luci-static/resources/view/network/interfaces.js
sed -i "s+s = m.section(form.GridSection, 'wifi-device', _('Wireless Overview'));+s = m.section(form.GridSection, 'wifi-device', _('Wireless Overview'), _('<b>SmoothWAN Note</b>: By default, the first <em>Generic</em> radio is the 5Ghz adapter, second is 2.4Ghz.'));+" /www/luci-static/resources/view/network/wireless.js

#Crude fixes
printf "config config 'Setup'\n option enabled '0'" > /etc/config/ntopconf
echo > /etc/hotplug.d/net/macs
echo -e 'nameserver ::1\nnameserver  127.0.0.1' > /etc/resolv.conf
/usr/bin/chattr +i /etc/resolv.conf &
/etc/init.d/apache2 disable &
cp /etc/factoryconfig/rt_tables /etc/iproute2/rt_tables &

echo "DISTRIB_DESCRIPTION='SmoothWAN v0.99.5HF5'" >> /etc/openwrt_release &
exit 0