# files/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.ipaddr='192.168.3.1'
commit network
set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-iface[0].ssid='SmoothWanSetup'
commit wireless
set system.@system[-1].hostname='SmoothWAN'
add system led
set system.@led[0]=led
set system.@led[0].name='Load'
set system.@led[0].sysfs='led0'
set system.@led[0].trigger='heartbeat'
add system led
set system.@led[1].name='Boot Done'
set system.@led[1].sysfs='led1'
set system.@led[1].trigger='none'
commit system
set dhcp.lan.dhcp_option='6,10.202.0.1,192.168.3.1'
add_list dhcp.cfg01411c.server='8.8.8.8'
commit dhcp
del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='Speedify'
commit firewall
set network.Speedify=interface
set network.Speedify.proto='none'
set network.Speedify.device='connectify0'
commit network
EOI

