# devconfigs/npir6s/files/etc/rc.local
#The following parameters forked to background instead of services for easy tweaking and diagnostics.

#Enable a scheduler for root disk for better LXC performance
echo "bfq" > /sys/class/block/mmcblk0/queue/scheduler &

#Run iperf3 server
iperf3 -s &

#P2P logging
sh /usr/lib/p2plog/p2p.sh &

#Load USB ethernet module late (NanoPi)
modprobe usbnet

exit 0
