# devconfigs/x86_64/files/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.ipaddr='172.17.17.2'
set network.AutoWANOne=interface
set network.AutoWANOne.proto='dhcp'
set network.AutoWANOne.device='eth1'
set network.AutoWANOne.metric='10'
set network.AutoWANTwo=interface
set network.AutoWANTwo.proto='dhcp'
set network.AutoWANTwo.device='eth2'
set network.AutoWANTwo.metric='20'
set network.AutoWANThree=interface
set network.AutoWANThree.proto='dhcp'
set network.AutoWANThree.device='eth3'
set network.AutoWANThree.metric='30'
set network.AutoWANFour=interface
set network.AutoWANFour.proto='dhcp'
set network.AutoWANFour.device='eth4'
set network.AutoWANFour.metric='40'
commit network
set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-iface[0].ssid='SmoothWAN Setup'
set wireless.@wifi-device[0].channel='auto'
set wireless.@wifi-device[0].hwmode='11g'
set wireless.@wifi-device[0].htmod='HT40'
set wireless.@wifi-device[0].country='US'
set wireless.@wifi-device[0].short_gi_40='0'
set wireless.@wifi-iface[0].encryption='psk2+ccmp'
set wireless.@wifi-iface[0].key='brassworld'
set wireless.@wifi-device[0].cell_density='0'
commit wireless
set system.@system[-1].hostname='SmoothWAN'
add system led
set system.@led[0]=led
set system.@led[0].name='System Load Indicator'
set system.@led[0].sysfs='led0'
set system.@led[0].trigger='heartbeat'
add system led
set system.@led[1].name='Boot Done'
set system.@led[1].sysfs='led1'
set system.@led[1].trigger='none'
commit system
set dhcp.cfg01411c.allservers='1'
set dhcp.cfg01411c.rebind_protection='0'
set dhcp.lan.force='1'
set dhcp.lan.force='1'
set dhcp.lan.ra='server'
set dhcp.lan.ra_default='2'
set dhcp.lan.ra_slaac='0'
del dhcp.lan.ra_flags
add_list dhcp.lan.ra_flags='none'
set dhcp.lan.dhcpv6='server'
commit dhcp
del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='Speedify'
add_list firewall.cfg03dc81.network='AutoWANOne'
add_list firewall.cfg03dc81.network='AutoWANTwo'
add_list firewall.cfg03dc81.network='AutoWANThree'
add_list firewall.cfg03dc81.network='AutoWANFour'
commit firewall
set network.Speedify=interface
set network.Speedify.proto='none'
set network.Speedify.device='connectify0'
set network.SPEEDIFY.proto='static'
set network.SPEEDIFY.ipaddr='10.202.0.2'
set network.SPEEDIFY.netmask='255.255.255.0'
set network.SPEEDIFY.gateway='10.202.0.1'
add_list network.SPEEDIFY.dns='10.202.0.1'
commit network
set nlbwmon.cfg015d22.database_directory='/usr/share/nlbwmon'
set nlbwmon.cfg015d22.commit_interval='10m'
commit nlbwmon
set system.ntp.use_dhcp='0'
del system.ntp.server
add_list system.ntp.server='162.159.200.123'
commit system
set luci.apply.holdoff='1'
commit luci

EOI

