# devconfigs/bpir3/files/etc/hotplug.d/iface/30-adghome
#!/bin/sh

/etc/init.d/adguardhome restart &

exit 0