# devconfigs/slateax/files/etc/uci-defaults/99-custom
uci -q batch << EOI
del network.wan
del network.wan6
set network.lan.ipaddr='172.17.17.2'
set network.PortOne=interface
set network.PortOne.proto='dhcp'
set network.PortOne.device='PortOne'
set network.PortOne.metric='10'
set network.PortTwo=interface
set network.PortTwo.proto='dhcp'
set network.PortTwo.device='PortTwo'
set network.PortTwo.metric='20'
set network.AutoUSBeth0=interface
set network.AutoUSBeth0.proto='dhcp'
set network.AutoUSBeth0.device='eth0'
set network.AutoUSBeth0.metric='21'
set network.AutoUSBeth1=interface
set network.AutoUSBeth1.proto='dhcp'
set network.AutoUSBeth1.device='eth1'
set network.AutoUSBeth1.metric='31'
set network.AutoUSB1=interface
set network.AutoUSB1.proto='dhcp'
set network.AutoUSB1.device='usb0'
set network.AutoUSB1.metric='30'
set network.AutoUSB2=interface
set network.AutoUSB2.proto='dhcp'
set network.AutoUSB2.device='usb1'
set network.AutoUSB2.metric='40'
del network.cfg030f15.ports
add_list network.cfg030f15.ports='eth2'
commit network

set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-iface[0].ssid='SmoothWAN Setup'
set wireless.@wifi-device[0].channel='auto'
set wireless.@wifi-device[0].channels='36 40 44 48'
set wireless.@wifi-device[0].band='5g'
set wireless.@wifi-device[0].htmod='HE80'
set wireless.@wifi-device[0].country='US'
set wireless.@wifi-iface[0].encryption='psk2+ccmp'
set wireless.@wifi-iface[0].key='brassworld'
commit wireless

set system.@system[-1].hostname='SmoothWAN'
set system.@system[-1].log_size='1024'
del system.led_wan
add system led
set system.@led[0].trigger='default-on'
set system.@led[0].sysfs='white_led'
set system.@led[0].name='Boot Done'
add system led
set system.@led[1].trigger='none'
set system.@led[1].sysfs='blue_led'
set system.@led[1].name='Boot Done'
commit system

set dhcp.cfg01411c.allservers='1'
set dhcp.cfg01411c.rebind_protection='0'
add_list dhcp.cfg01411c.server='127.0.0.1#5353'
set dhcp.lan.force='1'
del dhcp.lan.ra
del dhcp.lan.ra_default
del dhcp.lan.ra_slaac
del dhcp.lan.dhcpv6
del dhcp.lan.ra_flags
add_list dhcp.lan.ra_flags='none'
set dhcp.lan.dhcpv6='server'
set dhcp.@dnsmasq[0].dhcpscript='/usr/sbin/iptmon'
commit dhcp

del firewall.cfg03dc81.network
add_list firewall.cfg03dc81.network='Speedify'
add_list firewall.cfg03dc81.network='PortOne'
add_list firewall.cfg03dc81.network='PortTwo'
add_list firewall.cfg03dc81.network='AutoUSB1'
add_list firewall.cfg03dc81.network='AutoUSB2'
add_list firewall.cfg03dc81.network='AutoUSBeth0'
add_list firewall.cfg03dc81.network='AutoUSBeth1'
commit firewall

set network.Speedify=interface
set network.Speedify.proto='none'
set network.Speedify.device='connectify0'
set network.Speedify.peerdns='0'
add_list network.Speedify.dns='10.202.0.1'
commit network

set system.ntp.use_dhcp='0'
del system.ntp.server
add_list system.ntp.server='162.159.200.123'
commit system

set luci.apply.holdoff='1'
commit luci

add_list vpn-policy-routing.config.supported_interface='PortOne'
add_list vpn-policy-routing.config.supported_interface='PortTwo'
add_list vpn-policy-routing.config.supported_interface='AutoUSB1'
add_list vpn-policy-routing.config.supported_interface='AutoUSB2'
add_list vpn-policy-routing.config.supported_interface='AutoUSBeth0'
add_list vpn-policy-routing.config.supported_interface='AutoUSBeth1'
add_list vpn-policy-routing.config.supported_interface='Speedify'
commit vpn-policy-routing

set luci_statistics.collectd_sensors.enable='1'
set luci_statistics.collectd_sensors.IgnoreSelected='1'
del luci_statistics.cfg1b3a2e
set luci_statistics.collectd_iptables.enable='1'
set luci_statistics.cfg1a3a2e.name='LAN'
set luci_statistics.cfg1a3a2e.chain='zone_lan_postrouting'
set luci_statistics.collectd_ping.enable='1'
del luci_statistics.collectd_ping.Hosts
add_list luci_statistics.collectd_ping.Hosts='10.202.0.1'
add_list luci_statistics.collectd_ping.Hosts='10.202.0.5'
add_list luci_statistics.collectd_ping.Hosts='8.8.8.8'
set luci_statistics.collectd_ping.AddressFamily='any'
del luci_statistics.collectd_uptime.enable
del luci_statistics.collectd_network.enable
set luci_statistics.collectd.Include='/etc/collectd/conf.d/*.conf'
set luci_statistics.collectd_uptime.enable='1'
set luci_statistics.collectd_interface.IgnoreSelected='1'
add_list luci_statistics.collectd_interface.Interfaces='br-lan'
add_list luci_statistics.collectd_interface.Interfaces='ip6tnl0'
add_list luci_statistics.collectd_interface.Interfaces='bond0'
add_list luci_statistics.collectd_interface.Interfaces='lo'
add_list luci_statistics.collectd_interface.Interfaces='gretap0'
add_list luci_statistics.collectd_interface.Interfaces='ip6gre0'
add_list luci_statistics.collectd_interface.Interfaces='gre0'
set luci_statistics.collectd_rrdtool.DataDir='/etc/rrd'
commit luci_statistics
set ttyd.cfg01a8ea.ssl='1'
set ttyd.cfg01a8ea.ssl_cert='/etc/nginx/nginx.cer'
set ttyd.cfg01a8ea.ssl_key='/etc/nginx/nginx.key'
commit ttyd
add fstab mount # =cfg024d78
set fstab.@mount[-1].enabled='1'
set fstab.@mount[-1].target='/mnt/smoothwan'
set fstab.@mount[-1].uuid='337d0b9f-8716-42b2-b316-312811aef984'
commit fstab

EOI


#Replace diagnostics page
cp /etc/sysupgdiff/diagnostics.js /www/luci-static/resources/view/network/diagnostics.js
#Hints
sed -i "s+Choose the network(s) you want to attach to this wireless interface or fill out the <em>custom</em> field to define a new network.+Choose the network(s) you want to attach to this wireless interface or fill out the <em>custom</em> field to define a new network.<br><br><b>SmoothWAN Note!<br><br>If this is a Wireless WAN, head to Network->Interfaces and edit the automatically created interface to set the metric to 10.</b>+g" /www/luci-static/resources/view/network/wireless.js
sed -i "s+ _('Bring up on boot')+ _('Bring up on boot'), _('<br><br><b>SmoothWAN Note!<br><br>If this is a WAN connection, head to Advanced Settings tab above and set metric to 10')+" /www/luci-static/resources/view/network/interfaces.js
sed -i "s+s = m.section(form.GridSection, 'wifi-device', _('Wireless Overview'));+s = m.section(form.GridSection, 'wifi-device', _('Wireless Overview'), _('<b>SmoothWAN Note</b>: By default, the first <em>Generic</em> radio is the 5Ghz adapter, second is 2.4Ghz.'));+" /www/luci-static/resources/view/network/wireless.js